# server/dawnpilot.service
[Unit]
Description=DawnPilot weather cache service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dawnpilot
Group=dawnpilot
WorkingDirectory=/opt/dawnpilot
EnvironmentFile=/etc/dawnpilot.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/dawnpilot/dawnpilot_server.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/dawnpilot
MemoryMax=200M

[Install]
WantedBy=multi-user.target
